# debian/roush-simple.roush.upstart.in
# roush - a web api for cluster orchestration orchestration (aaS!)
#

description     "roush web api"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

script
    DAEMON_ARGS="-c /etc/roush/roush.conf"
    [ -e /etc/default/roush ] && . /etc/default/roush

    /usr/bin/roush ${DAEMON_ARGS}
end script
